# pi_web_app/ps4-monitor.service
[Unit]
Description=PS4 Controller Web Monitor
After=network-online.target bluetooth.target
Wants=network-online.target bluetooth.target

[Service]
Type=simple
User=chiciu
Group=bluetooth
WorkingDirectory=/home/chiciu/Documents/4_axis_robotic_arm/pi_web_app
ExecStart=/home/chiciu/Documents/4_axis_robotic_arm/pi_web_app/venv/bin/python /home/chiciu/Documents/4_axis_robotic_arm/pi_web_app/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/chiciu/Documents/4_axis_robotic_arm/pi_web_app/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=SDL_VIDEODRIVER=dummy
Environment=SDL_AUDIODRIVER=dummy

[Install]
WantedBy=multi-user.target